# systemd/voicein.service
[Unit]
Description=VoiceIn tray dictation (STT + LLM paste)

[Service]
Type=simple

Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=%t

WorkingDirectory=@VOICEIN_REPO@

ExecStart=@VOICEIN_VENV_PYTHON@ -m voicein run

Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
